# steam-limiter.service
[Unit]
Description=Steam playtime limiter
After=graphical-session.target

[Service]
Type=simple
ExecStart=%h/.local/bin/steam-limiter.sh
Restart=always
RestartSec=10
Environment=LIMIT_MINUTES=60
Environment=POLL_SECONDS=10

[Install]
WantedBy=default.target
